# supervisord.conf
[supervisord]

nodaemon=true               ; (start in foreground if true;default false)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
user=root                    ; (default is current user, required if root)

[program:app]
command=npm start --prefix /app
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

